# check-metadata-server.timer
[Unit]
Description=Check Metadata Server for Update

[Timer]
OnCalendar=hourly

[Install]
WantedBy=multi-user.target
